# Ubuntu Config Files/rc.local
#!/bin/bash -e
# /etc/rc.local

#enable wake on lan
sleep 5
ethtool -s enp3s0 wol g

exit 0

